# .config/autostart/evsieve.desktop
[Desktop Entry]
Exec=/home/iprf/scripts/evsieve
Icon=application-x-shellscript
Name=evsieve
Type=Application
X-KDE-AutostartScript=true
