# .config/autostart/org.kde.kclockd-autostart.desktop
[Desktop Entry]
DBusActivatable=true
Exec=kclockd
Name=org.kde.kclockd-autostart
Type=Application
X-Flatpak=org.kde.kclockd-autostart
