# .config/autostart/sway-audio-idle-inhibit.desktop
[Desktop Entry]
Exec=sway-audio-idle-inhibit
Icon=
Name=sway-audio-idle-inhibit
Path=
Terminal=False
Type=Application
