# rootlinks/udev/90-external-mouse.rules
SUBSYSTEM=="input", KERNEL=="mouse[0-9]*", RUN+="/usr/bin/su iprf -c /home/iprf/scripts/set-touchpad"
